# mysite.ini
[program:mysite]
user=root                                         ; 用哪个用户启动
command=celery worker --app=task -l info             ; 启动命令
stdout_logfile=/var/log/supervisor/celeryd_out.log   ; stdout 日志输出位置
stderr_logfile=/var/log/supervisor/celeryd_err.log   ; stderr 日志输出位置
stdout_logfile_maxbytes=20MB                         ; stdout 日志文件大小，默认 50MB
stdout_logfile_backups=20                            ; stdout 日志文件备份数
autostart=true                                       ; 在 supervisord 启动的时候自动启动
autorestart=true                                     ; 程序异常退出后自动重启
startsecs=10                                         ; 启动 10 秒后没有异常退出，就当作已经正常启动
environment=DJANGO_SECRET_KEY='3ce2*w)u5p(663(cvzpo@-$08+wvi$^^n4x#$@29@x^-nhi^qy'